# infra/systemd/flowbiz-beauty-api.staging.service.example
[Unit]
Description=FlowBiz Beauty API Service (staging)
After=network.target

[Service]
Type=simple
User=flowbiz
Group=flowbiz
WorkingDirectory=/opt/flowbiz/clients/flowbiz-client-beauty/current
Environment=NODE_ENV=production
EnvironmentFile=/etc/flowbiz/flowbiz-beauty-staging.env
ExecStart=/usr/bin/node apps/api/src/server.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
